# com.vibebot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  VibeBot launchd agent — runs the daily digest at 9:00 AM.

  Useful commands:
    launchctl start com.vibebot.daily        # trigger immediately
    launchctl stop com.vibebot.daily         # stop a running job
    launchctl unload ~/Library/LaunchAgents/com.vibebot.plist  # disable
    tail -f ~/Library/Logs/vibebot.log       # view logs
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.vibebot.daily</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>-m</string>
        <string>vibebot.main</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/utsavmehta/PersonalProjects/VibeBot</string>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/utsavmehta/Library/Logs/vibebot.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/utsavmehta/Library/Logs/vibebot.log</string>
</dict>
</plist>
