# launchd/com.vibebot.f1.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
    "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.vibebot.f1</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>-m</string>
        <string>vibebot.workers.f1_check</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/utsavmehta/PersonalProjects/VibeBot</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONPATH</key>
        <string>/Users/utsavmehta/PersonalProjects/VibeBot</string>
    </dict>

    <key>StartInterval</key>
    <integer>300</integer>

    <key>StandardOutPath</key>
    <string>/tmp/vibebot_f1.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/vibebot_f1_err.log</string>
</dict>
</plist>
