# launchd/com.vibebot.sports.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.vibebot.sports</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/node</string>
        <string>/Users/utsavmehta/PersonalProjects/VibeBot/dist/main.js</string>
        <string>--plugin</string>
        <string>sports</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/utsavmehta/PersonalProjects/VibeBot</string>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>10</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/utsavmehta/Library/Logs/vibebot-sports.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/utsavmehta/Library/Logs/vibebot-sports.log</string>
</dict>
</plist>
