# launchd/com.vibebot.mlx-audio.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.vibebot.mlx-audio</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/utsavmehta/.local/bin/mlx_audio.server</string>
        <string>--port</string>
        <string>8080</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/utsavmehta</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/utsavmehta/Library/Logs/vibebot-mlx-audio.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/utsavmehta/Library/Logs/vibebot-mlx-audio.log</string>
</dict>
</plist>
