# launchd/com.vibebot.india_news.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.vibebot.india_news</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/node</string>
        <string>/Users/utsavmehta/PersonalProjects/VibeBot/dist/main.js</string>
        <string>--plugin</string>
        <string>india_news</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/utsavmehta/PersonalProjects/VibeBot</string>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/utsavmehta/Library/Logs/vibebot-india_news.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/utsavmehta/Library/Logs/vibebot-india_news.log</string>
</dict>
</plist>
